# rootfs/etc/s6-overlay/s6-rc.d/webserver/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: WeeWX
# ==============================================================================

bashio::log.info "Starting web server..."
python -m http.server -d /config/public_html/ 80
