# rootfs/etc/s6-overlay/s6-rc.d/weewx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: WeeWX
# ==============================================================================

#

if ! bashio::fs.file_exists /config/weewx.conf;
then
    bashio::log.info "Creating default weewx config"
    weectl station create --no-prompt /config

    bashio::log.info "Installing weexw-interceptor"
    weectl extension install --config /config/weewx.conf --yes https://github.com/matthewwall/weewx-interceptor/archive/master.zip 
    
    bashio::log.info "Installing weexw-mqtt"
    weectl extension install --config /config/weewx.conf --yes https://github.com/matthewwall/weewx-mqtt/archive/master.zip
fi

bashio::log.info "Starting WeeWX..."
weewxd --config /config/weewx.conf
